# system_report.cron
*/5 * * * * root /usr/local/bin/system_report.sh
